# Crontab/pi@buoypi.cron
# Crontab for buoy pi

0 22,23 * * * /home/pi/profile_manager/utils/L4tides/run_L4buoy_tidal_model.sh >> /home/pi/profile_manager/utils/L4tides/Logs/L4slack.log 2>&1
